# docker/node-exporter/crontab
* * * * * /root/scripts/file.sh | sponge /data/node-exporter/file.prom
