# docker/node-exporter/crontab
* * * * * /root/scripts/file.sh | sponge $METRICS_ROOT/file.prom
