# docker/node-exporter/crontab
* * * * * /root/scripts/node_exporter.sh
